# systemd/v7-health.service
[Unit]
Description=V7 Health Loop
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'while true; do v7-egress-history; v7-egress-stability; v7-egress-load; v7-egress-diagnose; v7-state-merge; v7-user-desired-state-save; v7-state-json-save; v7-users-autoswitch; sleep 30; done'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
