# systemd/v7-path-sanity.service
[Unit]
Description=V7 client path sanity check
After=network-online.target v7-mss-clamp.service v7-direct-autosync.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/v7-path-sanity-check

